# etc/profile.d/welcome.sh
#!/bin/sh
# Part of cansel-config
#
# See LICENSE file for copyright and license details

# Should be installed to /etc/profile.d/welcome.sh to force welcome
# to run at initial login

welcome

